# console.upstart
description "Loads console font and keyboard table"

start on pld.sysinit-done

task
console output

pre-start script
	[ -f /etc/sysconfig/console ] || exit 0

	[ -f /bin/loadkeys -o -f /usr/bin/loadkeys ]
	[ -f /bin/setfont -o -f /usr/bin/setfont ]
end script

script
	export USE_UPSTART=no
	exec /etc/rc.d/init.d/console start >/dev/null
end script
